# roles/chirpstack_v3_as/files/chirpstack-application-server.service
[Unit]
Description=ChirpStack Application Server
Documentation=https://www.chirpstack.io/
Wants=network-online.target
After=network-online.target

[Service]
User=appserver
Group=appserver
ExecStart=/usr/local/bin/chirpstack-application-server
Restart=on-failure
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
Alias=lora-app-server.service
